# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ollama]
command=ollama serve
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama.err

[program:brain]
command=node /app/src/agi/brain.js
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/brain.log
stderr_logfile=/var/log/brain.err
depends_on=ollama

[program:dashboard]
command=node /app/src/ui/server.js
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/dashboard.log
stderr_logfile=/var/log/dashboard.err
depends_on=brain
