# supervisord.conf
[supervisord]
nodaemon=true

[program:engine]
command=/app/emerald-engine
autostart=true
autorestart=true
stdout_logfile=/tmp/engine.log
stderr_logfile=/tmp/engine.log
