# supervisord.conf
[supervisord]
nodaemon=true

[program:ollama]
command=ollama serve
autostart=true
autorestart=true

[program:node]
command=node src/agi/brain.js
autostart=true
autorestart=true
