# my-aladdin-system/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:aladdin-core]
command=python main.py
directory=/app
stdout_logfile=/data/logs/aladdin-core.log
stderr_logfile=/data/logs/aladdin-core.err
autostart=true
autorestart=true
startretries=10
stopwaitsecs=10

[program:micro-earning]
command=python micro_earning.py
directory=/app
stdout_logfile=/data/logs/micro-earning.log
stderr_logfile=/data/logs/micro-earning.err
autostart=true
autorestart=true
startretries=10
stopwaitsecs=5

[program:ollama]
command=ollama serve
directory=/app
stdout_logfile=/data/logs/ollama.log
stderr_logfile=/data/logs/ollama.err
autostart=true
autorestart=true
startretries=10
stopwaitsecs=30

[program:cloudflared]
command=bash /app/scripts/cloudflare_tunnel.sh
directory=/app
stdout_logfile=/data/logs/cloudflared.log
stderr_logfile=/data/logs/cloudflared.err
autostart=false
autorestart=true
startretries=5
stopwaitsecs=5
